# common/.config/systemd/user/cleanup-files.timer
[Unit]
Description=Daily clean-up of the Downloads folder

[Timer]
OnStartupSec=1min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target

